# assets/systemd/servika-db-backup.service
[Unit]
Description=Servika panel database backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/servika-db-backup
